# roles/rtsp_simple_server/files/rtsp-simple-server.service
[Unit]
Wants=network.target

[Service]
ExecStart=/usr/bin/rtsp-simple-server /etc/rtsp-simple-server.yaml

[Install]
WantedBy=multi-user.target