# roles/prometheus_exporters/files/node-exporter.service
[Unit]
Description=Node-exporter
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/tmp/
ExecStart=/usr/bin/node_exporter
Restart=always

[Install]
WantedBy=multi-user.target
